# dmsgpty-tcp@.service
[Unit]
Description=Skywire dmsgpty-host session (%i)
After=network-online.target

[Service]
# Single-conn mode: dmsgpty-host reads FD 0 as net.Conn, runs one
# noise handshake + authorize + session pipeline, exits when the
# session closes. systemd hands the accepted socket on stdin
# because the companion .socket sets Accept=yes.
#
# --no-dmsg keeps this daemon off the dmsg overlay so it can safely
# share the visor's identity via --sk-from-visor without contending
# for the visor's dmsg-discovery entry.
#
# Paths match the skywire-bin package layout. Edit if your install
# uses non-default locations:
#   /opt/skywire/skywire.json  — visor config; --sk-from-visor reads
#                                 its .sk field for the host identity
#   /opt/skywire/dmsgpty.json  — host whitelist + daemon settings;
#                                 auto-generated below if missing
#
# The unified skywire binary supersedes the separate dmsgpty-host
# binary that used to live at /usr/bin/dmsgpty-host. `skywire dmsg
# pty host` carries the same flag surface (--no-dmsg,
# --sk-from-visor, --confpath, --listen-fd) so socket-activation
# semantics are unchanged.
ExecStartPre=/bin/bash -c '[[ -f /opt/skywire/dmsgpty.json ]] || /usr/bin/skywire dmsg pty host confgen /opt/skywire/dmsgpty.json'
ExecStart=/usr/bin/skywire dmsg pty host \
    --no-dmsg \
    --sk-from-visor /opt/skywire/skywire.json \
    --confpath /opt/skywire/dmsgpty.json \
    --listen-fd 0
StandardInput=socket
StandardOutput=journal
StandardError=journal

# Per-session process: no Restart= — the session is over when the
# user logs out, just like sshd's per-session children.
Type=simple
