# skywire.service
[Unit]
Description=Skywire
Wants=network-online.target

[Service]
Environment="GODEBUG=madvdontneed=1"
Type=simple
Environment="SYSTEMDCHILD=1"
ExecStart=/usr/bin/skywire visor -ps -q http --loglvl debug
Restart=always
RestartSec=20
TimeoutSec=30
# Hard memory limit — systemd kills the visor cleanly if exceeded.
# Prevents OOM killer from randomly picking victims on low-RAM boards.
MemoryMax=80%

[Install]
WantedBy=multi-user.target
