# skywire-bin/skywire-docker-update.service
[Unit]
Description=Skywire deployment docker auto-updater
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/skywire-docker-update
